# systemd/ads-b-graphs-7d.service
[Unit]
Description=ADS-B Graphs 7d

[Service]
Type=oneshot
ExecStart=/usr/bin/make-collectd-graphs.sh 7d
