# systemd/ads-b-graphs-24h.service
[Unit]
Description=ADS-B Graphs 24h

[Service]
Type=oneshot
ExecStart=/usr/bin/make-collectd-graphs.sh 24h
