# systemd/ads-b-graphs-1h.timer
[Unit]
Description=ADS-B Graphs 1h

[Timer]
OnBootSec=10
OnCalendar=*:0/10
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
