# systemd/ads-b-graphs-7d.timer
[Unit]
Description=ADS-B Graphs 7d

[Timer]
OnBootSec=10
OnCalendar=*-*-* *:00:00
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
