# systemd/ads-b-graphs-365d.service
[Unit]
Description=ADS-B Graphs 365d

[Service]
Type=oneshot
ExecStart=/usr/bin/make-collectd-graphs.sh 365d
