# systemd/ads-b-graphs-6h.timer
[Unit]
Description=ADS-B Graphs 6h

[Timer]
OnBootSec=10
OnCalendar=*:0/15
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
