# systemd/ads-b-graphs-30d.service
[Unit]
Description=ADS-B Graphs 30d

[Service]
Type=oneshot
ExecStart=/usr/bin/make-collectd-graphs.sh 30d
